# ssh.conf
[program:ssh]
autorestart=true
command=/root/start_ssh.sh
stderr_logfile = /dev/stdout ;/var/log/supervisord/sshd-stderr.log
stdout_logfile = /dev/stdout ;/var/log/supervisord/sshd-stdout.log
